# bdre-scripts/local/bin/bdre
#!/bin/bash
#
# chkconfig: 35 90 12
# description: Foo server
#
# Get function from functions library
. /etc/init.d/functions
# Start the service BDRE
USER=%USER%
BDRE_HOME=~$USER/bdre
BDRE_APPS_HOME=~$USER/bdre_apps
LOCKFILE=/var/lock/subsys/bdre
DAEMON=$BDRE_HOME/bdre-scripts/execution/run-ui.sh # Introduce the bdre's location here
LOGDIR=/var/log/bdre  # Log directory to use
PIDFILE=/var/run/bdre/supervisor.pid

DAEMON_OPTS="-p $PIDFILE $LOGDIR"
BDRE_SHUTDOWN_TIMEOUT=15


start() {
        echo -n "Starting bdre: "
        for dir in $(dirname $PIDFILE)
        do
            mkdir -p $dir
            chown -R $USER $dir
        done
	# Check if already running
        if [ -e $PIDFILE ] && checkpid $(cat $PIDFILE) ; then
            echo "already running"
            return 0
        fi
	# the supervisor itself will setuid down to $USER
        sudo /sbin/runuser -s /bin/bash $USER -c "$DAEMON $DAEMON_OPTS/service.log"
	    ret=$?
        base=$(basename $0)
        if [ $ret -eq 0 ]; then
            sleep 5
            test -e $PIDFILE && checkpid $(cat $PIDFILE)
            ret=$?
        fi
        if [ $ret -eq 0 ]; then
            touch $LOCKFILE
            success $"$base startup"
        else
            failure $"$base startup"
        fi
        echo
        return $ret

}
# Restart the service bdre
stop() {
	 if [ ! -e $PIDFILE ]; then
            success "bdre is not running"
            return 0
        fi

        echo -n "Shutting down bdre: "

        BDRE_PID=`cat $PIDFILE 2>/dev/null`
        if [ -n "$BDRE_PID" ]; then
          kill -TERM ${BDRE_PID} &>/dev/null
          for i in `seq 1 ${BDRE_SHUTDOWN_TIMEOUT}` ; do
            kill -0 ${BDRE_PID} &>/dev/null || break
            sleep 1
          done
          kill -KILL ${BDRE_PID} &>/dev/null
        fi
        echo
        rm -f $LOCKFILE $PIDFILE
        return 0
}

restart() {
  stop
  start
}

### main logic ###
case "$1" in
  start)
        start
        ;;
  stop)
        stop
        ;;
  status)
        status -p $PIDFILE supervisor
        ;;
  restart|reload)
        restart
        ;;
  condrestart)
        [ -f $LOCKFILE ] && restart || :
        ;;
    *)
        echo "Usage: bdre {start|stop|status|reload|restart|condrestart"
        exit 1
        ;;

esac
exit $?